# bin/crontab
#!/usr/bin/env php
<?php

require dirname(__FILE__) . '/runtime.php';