# bin/crontab
#!/usr/bin/env php
<?php

error_reporting(E_ALL & ~E_NOTICE);
if (PHP_SAPI !== 'cli') exit('CronTab only run under cli mode.');

define('CRON_DIR', dirname(__FILE__));

// Change current work dir
chdir(CRON_DIR);

// Cli arguments init.
$executor = __FILE__;

// Master trip
$cron_master = false;
if ($index = array_search('--master', $argv)) {
    unset($argv[$index]);
    $argv = array_values($argv);
    $cron_master = true;
}
unset($argv[0]);
$commands = array_values($argv);

// Init cron
require dirname(CRON_DIR) . '/src/AutoLoader.php';
AutoLoader::register();
$ini = new \CronTab\IniParser(CRON_DIR . '/crontab.ini');
$cron = new \CronTab\CronTab($ini->get(), CRON_DIR);
$cron->registerExecutor($executor);

// Work mode
if ($commands) {
    if (count($commands) > 1) {
        foreach ($commands as $command) {
            $cron->dispatch($command);
        }
        exit;
    }

    $cron->execute($commands[0]);
    exit;
}

if ($cron_master) {
    // Pulse trigger
    if ((int)shell_exec('ps -ef | grep "crontab --master" | grep -v grep | wc -l') > 1) {
        exit;
    }

    // Start cron
    $cron->start();
} else {
    die ('[Usage]: use "' . $executor . ' --master".' . PHP_EOL);
}
