# bin/crontab
#!/usr/bin/env php
<?php

// load CronTab
$loader = include(__DIR__ . '/../src/bootstrap.php');

require dirname(__FILE__) . '/runtime.php';